# systemd/user@.service.d/mammutfs-home.conf
[Unit]
Requires=mammutfs-home@%i.service
